# udmabuf/debian/postinst
#!/bin/sh -e

depmod

systemctl enable  udmabuf.service
systemctl start   udmabuf.service
