# udmabuf/debian/prerm
#!/bin/sh -e

systemctl stop    udmabuf.service
systemctl disable udmabuf.service

