# fclkcfg/debian/postinst
#!/bin/sh -e

depmod

systemctl enable  fpga-clock.service
systemctl start   fpga-clock.service
