# fclkcfg/debian/prerm
#!/bin/sh -e

systemctl stop    fpga-clock.service
systemctl disable fpga-clock.service

